# infrastructure/launchagents/com.zenithjoy.topic-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zenithjoy.topic-worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/administrator/perfect21/zenithjoy/services/creator/scripts/topic-worker.py</string>
        <string>--apply</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/administrator/perfect21/zenithjoy</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/topic-worker.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/topic-worker.log</string>
</dict>
</plist>
