# infrastructure/launchagents/com.zenithjoy.pipeline-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zenithjoy.pipeline-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/administrator/perfect21/zenithjoy/services/creator/.venv/bin/python3</string>
        <string>/Users/administrator/perfect21/zenithjoy/services/creator/pipeline_worker/worker.py</string>
        <string>--apply</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/administrator/perfect21/zenithjoy</string>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/pipeline-worker.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/pipeline-worker.log</string>

    <!--
        PR-c：worker 全 HTTP 化，不再直连 SQLite。
        部署步骤：
          1. cp 本文件到 ~/Library/LaunchAgents/
          2. /usr/libexec/PlistBuddy -c "Set :EnvironmentVariables:ZENITHJOY_INTERNAL_TOKEN $(grep -m1 '^ZENITHJOY_INTERNAL_TOKEN=' ~/.credentials/zenithjoy-internal-token | cut -d= -f2-)" ~/Library/LaunchAgents/com.zenithjoy.pipeline-worker.plist
          3. launchctl unload ~/Library/LaunchAgents/com.zenithjoy.pipeline-worker.plist
          4. launchctl load ~/Library/LaunchAgents/com.zenithjoy.pipeline-worker.plist
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>APPS_API_BASE</key>
        <string>http://localhost:5200</string>
        <key>ZENITHJOY_INTERNAL_TOKEN</key>
        <string>PLACEHOLDER_TOKEN</string>
        <key>BRAIN_URL</key>
        <string>http://localhost:5221</string>
        <key>CREATOR_DEFAULT_NOTEBOOK_ID</key>
        <string>1d928181-4462-47d4-b4c0-89d3696344ab</string>
        <key>NAS_USER</key>
        <string>徐啸</string>
        <key>NAS_HOST</key>
        <string>100.110.241.76</string>
        <key>NAS_BASE</key>
        <string>/volume1/workspace/vault/zenithjoy-creator/content</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>LANG</key>
        <string>zh_CN.UTF-8</string>
        <key>TZ</key>
        <string>Asia/Shanghai</string>
    </dict>
</dict>
</plist>
