# infrastructure/launchagents/com.zenithjoy.topic-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zenithjoy.topic-worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/administrator/perfect21/zenithjoy/services/creator/scripts/topic-worker.py</string>
        <string>--apply</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/administrator/perfect21/zenithjoy</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/topic-worker.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/topic-worker.log</string>

    <!--
        PR-c：worker 全 HTTP 化，需要 apps/api base + 内部 token。
        部署步骤：
          1. cp 本文件到 ~/Library/LaunchAgents/
          2. /usr/libexec/PlistBuddy -c "Set :EnvironmentVariables:ZENITHJOY_INTERNAL_TOKEN $(grep -m1 '^ZENITHJOY_INTERNAL_TOKEN=' ~/.credentials/zenithjoy-internal-token | cut -d= -f2-)" ~/Library/LaunchAgents/com.zenithjoy.topic-worker.plist
          3. launchctl unload ~/Library/LaunchAgents/com.zenithjoy.topic-worker.plist
          4. launchctl load ~/Library/LaunchAgents/com.zenithjoy.topic-worker.plist
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>APPS_API_BASE</key>
        <string>http://localhost:5200</string>
        <key>ZENITHJOY_INTERNAL_TOKEN</key>
        <string>PLACEHOLDER_TOKEN</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>LANG</key>
        <string>zh_CN.UTF-8</string>
        <key>TZ</key>
        <string>Asia/Shanghai</string>
    </dict>
</dict>
</plist>
